# rootfs/etc/s6-overlay/s6-rc.d/cups/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting CUPS daemon..."

# Wait for D-Bus
while [ ! -e /run/dbus/system_bus_socket ]; do
    bashio::log.info "Waiting for D-Bus..."
    sleep 1
done

exec cupsd -f
