# rootfs/etc/s6-overlay/s6-rc.d/relayprint/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
set -euo pipefail

bashio::log.info "Starting RelayPrint API server on port 7779..."

# Change to the directory with the Python modules
cd /usr/local/bin

# Start the Print API server
exec python3 -c "from print_api import app; app.run(host='0.0.0.0', port=7779, threaded=True)"
