# rootfs/etc/s6-overlay/s6-rc.d/dbus/run
#!/command/execlineb -P
with-contenv
s6-setuidgid root
exec dbus-daemon --system --nofork
